# runner/autosearch-runner.service
[Unit]
Description=Autosearch DPDK Runner Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dpdk
Group=dpdk
Environment=AUTOSEARCH_CONFIG=/etc/autosearch/runner.toml
ExecStart=/usr/local/bin/autosearch-runner
Restart=on-failure
RestartSec=30
TimeoutStopSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autosearch-runner

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/autosearch /tmp

[Install]
WantedBy=multi-user.target
